# scripts/target_system/drive_module_rpi/lib/systemd/system/robotdrivehost.service
[Unit]
Description=Drive Host for robot

[Service]
WorkingDirectory=/home/robot/
ExecStart=/opt/cotsbot/rundrivehost.sh
StanddardOutput=syslog
StandardError=syslog

[Install]
WantedBy=multi-user.target
