# scripts/target_system/drive_module_rpi/lib/systemd/system/robotmaincamera.service
[Unit]
Description=Webcam broadcaster for robot

[Service]
WorkingDirectory=/opt/cotsbot/
ExecStart=/opt/cotsbot/runrtspserver.sh
StandardOutput=syslog
StandardError=syslog

[Install]
WantedBy=multi-user.target
