# scripts/target_system/drive_module_rpi/lib/systemd/system/robotpwmdrivehost.service
[Unit]
Description=PWM Drive Host for robot

[Service]
WorkingDirectory=/home/robot/
ExecStart=/opt/cotsbot/runpwmdrivehost.sh
StanddardOutput=syslog
StandardError=syslog

[Install]
WantedBy=multi-user.target
